# common-files/nightshield-boot.service
[Unit]
Description=NightShield Boot Lock
DefaultDependencies=no
Before=sysinit.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nightshield-boot.sh

[Install]
WantedBy=basic.target